# ansible/proxmox-auto-dhcp/files/auto-dhcp.service
[Unit]
Description=Service for Auto Connect to new interfaces

[Service]
ExecStart=/usr/local/bin/auto-dhcp.sh
